# faas/crontab
# Ejecutar script cada hora
0 * * * * root node /app/index.js >> /var/log/cron.log 2>&1
